# deploy/node/aeronyx-server.service
# ============================================
# File: deploy/node/aeronyx-server.service
# ============================================
# Creation Reason:
# - Provide a repository-owned systemd unit template for production AeroNyx
#   privacy nodes.
#
# Modification Reason:
# - Add production systemd hardening while preserving TUN/VPN runtime access.
# - Make restart rate limiting explicit for production incident behavior.
# - Keep the configuration directory private because it can contain node
#   identity material, API credentials, and registration settings.
#
# Main Functionality:
# - Runs the Rust aeronyx-server binary under systemd.
# - Loads optional environment overrides from /etc/aeronyx/aeronyx.env.
# - Uses /etc/aeronyx/server.toml as the stable config path.
# - Restricts process privileges with capability and filesystem boundaries.
# - Lets systemd create the state/config/log directories before namespace
#   hardening is applied.
# - Pins `/etc/aeronyx` to mode 0700 so existing secure installations do not
#   produce a systemd directory-mode mismatch warning on every restart.
# - Uses explicit restart limits so transient failures can self-heal without
#   uncontrolled restart storms.
#
# Dependencies:
# - Installed by deploy/node/install.sh into /etc/systemd/system.
# - Reloaded/restarted by deploy/node/upgrade.sh.
# - Observed by deploy/node/healthcheck.sh and nodeboard restart controls.
#
# Main Logical Flow:
# 1. systemd waits for network-online.target.
# 2. aeronyx-server starts from @REPO_DIR@/target/release.
# 3. systemd applies runtime hardening that still permits TUN/VPN networking.
# 4. systemd restarts the service on process failure.
# 5. systemd applies explicit start-rate limits to avoid uncontrolled loops.
#
# Important Note for Next Developer:
# - Keep ExecStart compatible with crates/aeronyx-server/src/main.rs.
# - Do not add user traffic logging or destination inspection here.
# - install.sh replaces @REPO_DIR@ and @CONFIG_FILE@ before installation.
# - Do not enable PrivateDevices or ProtectHome here; production nodes need
#   /dev/net/tun and many current deployments place the repository under /root.
# - Do not relax ConfigurationDirectoryMode without a secrets migration and
#   an explicit review of every file written below `/etc/aeronyx`.
#
# Last Modified:
# v1.2.1-node-deploy - Pin private configuration directory mode to 0700.
# v1.2.0-node-deploy - Added explicit start-rate limits and start timeout.
# v1.1.1-node-deploy - Added systemd-managed state/config/log directories so
#                      ReadWritePaths cannot fail on missing paths.
# v1.1.0-node-deploy - Added production-safe systemd hardening.
# v1.0.0-node-deploy - Added production systemd unit template.
# ============================================

[Unit]
Description=AeroNyx Rust privacy node
After=network-online.target
Wants=network-online.target
StartLimitIntervalSec=300
StartLimitBurst=10

[Service]
Type=simple
User=root
WorkingDirectory=@REPO_DIR@
EnvironmentFile=-/etc/aeronyx/aeronyx.env
ExecStart=@REPO_DIR@/target/release/aeronyx-server start -c @CONFIG_FILE@
Restart=on-failure
RestartSec=5
KillSignal=SIGTERM
TimeoutStartSec=60
TimeoutStopSec=30
LimitNOFILE=1048576
ConfigurationDirectory=aeronyx
ConfigurationDirectoryMode=0700
StateDirectory=aeronyx
LogsDirectory=aeronyx

# Runtime hardening:
# - Keep CAP_NET_ADMIN/CAP_NET_RAW for TUN and VPN packet handling.
# - Keep CAP_NET_BIND_SERVICE for future privileged bind compatibility.
# - Keep /etc/aeronyx writable because first startup may persist generated
#   secrets and registration-related runtime state.
# - Do not use PrivateDevices=true; it would hide /dev/net/tun.
CapabilityBoundingSet=CAP_NET_ADMIN CAP_NET_BIND_SERVICE CAP_NET_RAW
AmbientCapabilities=CAP_NET_ADMIN CAP_NET_BIND_SERVICE CAP_NET_RAW
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=full
ReadWritePaths=/etc/aeronyx /var/lib/aeronyx /var/log/aeronyx
ProtectControlGroups=true
ProtectKernelTunables=true
ProtectKernelModules=true
RestrictSUIDSGID=true
LockPersonality=true
SystemCallArchitectures=native

[Install]
WantedBy=multi-user.target
